# supervisord.conf
[program:sd-agent]
command = /usr/bin/sd-agent/agent.py start
autorestart = true